# cmdr_recommend/cmdr_recommend.service
[Service]
ExecStart=/opt/cmdr/cmdr_recommend/run
Restart=always
StartLimitInterval=0
[Install]
WantedBy=default.target